# web2py/applications/spendtime/cron/crontab
01  *  *  *  *  root *applications/spendtime/controller/check.pyq