# web2py/applications/spendtime/cron/crontab
05   *  *  *  *  root *applications/spendtime/controller/check.py